# devise-agent/installers/macos/com.devise.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.devise.agent</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/devise-agent.app/Contents/MacOS/devise-agent</string>
        <string>--config</string>
        <string>/Library/Application Support/Devise/config.json</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/var/log/devise-agent.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/devise-agent-error.log</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>60</integer>
    
    <key>HangTimeout</key>
    <integer>120</integer>
</dict>
</plist>
